# crontab
# 定期支払い自動実行 - 毎日0時に実行
0 0 * * * root cd /code && python manage.py execute_recurring_payments >> /var/log/cron.log 2>&1

# セキュリティログ監視 - 1分ごとに実行
* * * * * root cd /code && python manage.py check_security_log >> /var/log/cron.log 2>&1

# デバッグログ監視 - 5分ごとに実行
*/5 * * * * root cd /code && python manage.py check_debug_log >> /var/log/cron.log 2>&1

# 空行が必要（cron仕様）
